# deploy/gunicorn.service

[Unit]
Description=Gunicorn for Flask app
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/srv/flaskapp
Environment="PATH=/srv/flaskapp/venv/bin"
ExecStart=/srv/flaskapp/venv/bin/gunicorn --workers 3 --bind unix:/run/flaskapp.sock wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
